# provisioning/mail-sandbox/general/etc/supervisor.d/dovecot.conf
[group:dovecot]
programs=dovecotd
priority=20

[program:dovecotd]
command = /opt/docker/bin/service.d/dovecot.sh
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
