# docker/php-nginx-dev/alpine-3/conf/etc/supervisor.d/nginx.conf
[group:nginx]
programs=nginxd
priority=20

[program:nginxd]
command = /opt/docker/bin/service.d/nginx.sh
process_name=%(program_name)s
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
