# mail-sandbox/conf/etc/supervisor.d/dovecot.conf
[group:dovecot]
programs=dovecot-daemon
priority=20

[program:dovecot-daemon]
command = /usr/sbin/dovecot -F
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
