# docker/php-official/7.4-alpine/conf/etc/supervisor.d/postfix.conf
[group:postfix]
programs=postfixd
priority=30

[program:postfixd]
directory = /etc/postfix
command = /opt/docker/bin/service.d/postfix.sh
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
