# _provisioning/php/ubuntu-12.04/etc/supervisor.d/php-fpm.conf
[program:php-fpm]
command = /usr/sbin/php5-fpm
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
