# base/centos-7/etc/supervisor.d/postfix.conf
[program:postfix]
directory = /etc/postfix
command = /opt/docker/bin/service.d/postfix.sh
startsecs = 0
autostart = false
autorestart = false