# provisioning/base/centos/etc/supervisor.d/cron.conf
[group:cron]
programs=crond
priority=30

[program:crond]
command = /sbin/crond -n
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
