# _provisioning/base/centos/etc/supervisor.d/cron.conf
[program:cron]
command = /sbin/crond -n
priority = 20
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
