# _provisioning/base/centos/etc/supervisor.d/cron.conf
[group:cron]
programs=cron-daemon
priority=20

[program:cron-daemon]
command = /sbin/crond -n
process_name=%(program_name)s
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
