# base/ubuntu-12.04/etc/supervisor.d/postfix.conf
[program:postfix]
directory = /etc/postfix
command = /opt/docker/bin/service.d/postfix.sh
priority = 20
startsecs = 0
autostart = false
autorestart = false
