# piwik/ubuntu-14.04/conf/etc/cron/crontab
5 * * * * application /usr/bin/php /application/code/piwik/console core:archive --url="$PIWIK_URL" > /dev/null
