# docker/piwik/latest/conf/etc/cron/crontab
5 * * * * application bash /opt/docker/bin/piwik-cron.sh &> /dev/null
