# mail-sandbox/ubuntu-14.04/conf/etc/supervisor.d/dovecot.conf
[group:dovecot]
programs=dovecotd
priority=20

[program:dovecotd]
command = /usr/sbin/dovecot -F
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
