# docker/php-dev/7.1-alpine/conf/etc/supervisor.d/blackfire-agent.conf
[group:blackfire-agent]
programs=blackfire-agentd
priority=25

[program:blackfire-agentd]
command = /opt/docker/bin/service.d/blackfire-agent.sh
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
