# php/debian-8/conf/etc/supervisor.d/log-php-error.conf
[program:log-php-error]
command = bash /opt/docker/bin/logwatch.sh php:error /var/log/php.error.log
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
