# php/ubuntu-15.04/conf/etc/supervisor.d/php-fpm.conf
[group:php-fpm]
programs=php-fpmd,php-fpm-log-fpm,php-fpm-log-slow,php-fpm-log-error,php-fpm-log-access
priority=20

[program:php-fpmd]
command = /usr/sbin/php5-fpm --nodaemonize
process_name=%(program_name)s
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm-log-fpm]
command = bash /opt/docker/bin/logwatch.sh php:fpm /var/log/php5-fpm/fpm.log
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm-log-slow]
command = bash /opt/docker/bin/logwatch.sh php:slow /var/log/php5-fpm/slow.log
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm-log-error]
command = bash /opt/docker/bin/logwatch.sh php:error /var/log/php5-fpm/error.log
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm-log-access]
command = bash /opt/docker/bin/logwatch.sh php:access /var/log/php5-fpm/access.log
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
